# .crontab
# m h  dom mon dow   command
25 12 * * *	(cd /home/miura/.Trash && /usr/bin/find . -mtime +7 -delete)
# 25 12 * * *	/usr/bin/find /home/miura/.Trash -mtime +7 -exec /bin/rm -fr {} \;
26 12 * * *	/usr/bin/trash-empty 7

# eof
